# systemd/feed.timer
[Unit]
Description=Timer for feed updates

[Timer]
Unit=feed.service
OnUnitInactiveSec=3h

[Install]
WantedBy=timers.target
