# systemd/feed.timer
[Unit]
Description=Timer for feed updates

[Timer]
Unit=feed.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
